# archive-deploy/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /app/archive-deploy/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
